# crontab
*/15 * * * * /etc/cron.d/chk_passwd_file.sh
